# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
# Wait for modules-load.d to fire (loads GPU modules from gpu.conf)
# and for udev to settle so DRM device nodes (/dev/dri/card*) are ready
# before X/Wayland tries to grab them. Without this, LightDM can race
# the kernel and fail with "no screens found" on slow-init GPUs.
#
# Additionally wait for ai-hw-detect.service so that the nouveau-legacy
# quirk (noaccel=1 on pre-Kepler cards) is applied BEFORE the X server
# tries to acceleration-probe the GPU.  Without this ordering, on a GT218
# host the X DDX driver opens glamor, hits the tiler bug, and freezes.
After=systemd-modules-load.service systemd-udev-settle.service ai-hw-detect.service
Wants=systemd-modules-load.service ai-hw-detect.service
