# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
# Wait for a DRI device before starting Xorg.
# systemd-modules-load runs in sysinit, so GPU modules should be loaded
# by the time graphical.target is reached, but just in case:
After=systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
# If /dev/dri/card0 still doesn't exist, try loading GPU modules manually.
# This covers edge cases where modules-load.d didn't fire (e.g. live ISO
# squashfs overlay timing). Failures are non-fatal (|| true).
ExecStartPre=/bin/sh -c 'if [ ! -e /dev/dri/card0 ]; then modprobe nouveau 2>/dev/null || true; modprobe amdgpu 2>/dev/null || true; sleep 1; fi'
