# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
# Wait for DKMS first-boot (builds nvidia module) and module loading
After=systemd-modules-load.service dkms-first-boot.service
Wants=systemd-modules-load.service

[Service]
# If /dev/dri/card0 doesn't exist yet, try loading GPU modules manually.
# Try nvidia first (DKMS-built), then nouveau, then amdgpu.
# If ALL fail, X will fall back to modesetting/fbdev/swrast.
ExecStartPre=/bin/sh -c 'if [ ! -e /dev/dri/card0 ]; then modprobe nvidia_drm modeset=1 2>/dev/null || modprobe nvidia 2>/dev/null || modprobe nouveau 2>/dev/null || modprobe amdgpu 2>/dev/null || true; sleep 2; fi'
# Restart up to 3 times if X fails (GPU module may load late)
Restart=on-failure
RestartSec=3
StartLimitBurst=3
StartLimitIntervalSec=30
