# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
After=systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
# Try loading GPU modules if /dev/dri/card0 isn't available yet.
# timeout ensures this never blocks boot for more than 5 seconds.
ExecStartPre=-/bin/sh -c 'timeout 5 sh -c "modprobe i915 2>/dev/null; modprobe nouveau 2>/dev/null; modprobe amdgpu 2>/dev/null" || true'
TimeoutStartSec=30
