# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
# Wait for modules-load.d to fire (loads GPU modules from gpu.conf)
After=systemd-modules-load.service
Wants=systemd-modules-load.service
