# profile/airootfs/etc/systemd/system/lightdm.service.d/ensure-gpu.conf
[Unit]
After=systemd-modules-load.service
Wants=systemd-modules-load.service

[Service]
# If /dev/dri/card0 doesn't exist, try loading GPU modules.
# Try nouveau FIRST (works on ALL nvidia GPUs, old and new).
# Only try nvidia_drm if nouveau didn't create a DRI device
# (nvidia-open only supports Maxwell+/SM5.0+).
# Modesetting X driver is the final fallback (10-modesetting-fallback.conf).
ExecStartPre=-/bin/sh -c 'if [ ! -e /dev/dri/card0 ]; then modprobe nouveau 2>/dev/null || modprobe amdgpu 2>/dev/null || modprobe nvidia_drm 2>/dev/null || modprobe i915 2>/dev/null || true; sleep 1; fi'
