# profile/airootfs/etc/systemd/system/lightdm.service.d/fallback.conf
# Session 44: wire ai-fbcon-fallback.service as lightdm's failure trigger.
#
# ai-fbcon-fallback.service is a tty2 emergency shell that should only
# activate when lightdm has genuinely given up (StartLimit burst hit,
# config broken, GPU missing post-detect).  Without this OnFailure=
# binding the fallback never fires and a broken display-manager leaves
# the user staring at a blank framebuffer with no recovery path.
#
# Intentionally NOT using Requires=/Requisite=: we want the fallback
# only on *failure*, not on clean lightdm shutdown (systemctl stop).
#
# Condition on !/run/archiso is already enforced inside ai-fbcon-fallback
# itself, so installing this drop-in on the live ISO is a no-op at
# activation time.
[Unit]
OnFailure=ai-fbcon-fallback.service
