# profile/airootfs/etc/systemd/system/lightdm.service.d/runtime-dir.conf
#
# Force systemd to create /run/lightdm before LightDM starts.
#
# The lightdm package ships a tmpfiles.d entry for /run/lightdm, but
# tmpfiles-setup runs early at boot and a systemctl restart of lightdm
# later will NOT recreate the dir if someone rm -rf'd /run/lightdm
# (the classic "I cleaned up /run" operator footgun).  RuntimeDirectory=
# makes systemd recreate it on every unit start, with the right owner.
#
# We also make sure LightDM starts only after user-sessions are allowed
# (otherwise PAM's pam_nologin denies 'arch' and autologin silently
# exits), after Plymouth has given up the framebuffer, and after the
# system-local-fs target so /etc is mounted rw (LightDM writes its
# .Xauthority files into /var/lib/lightdm which must be writable).
#
[Unit]
After=plymouth-quit-wait.service systemd-user-sessions.service local-fs.target
Wants=systemd-user-sessions.service
# If Plymouth is masked or disabled this After is a no-op, which is
# exactly what we want.
# StartLimit* must live in [Unit] since systemd 230; [Service] placement
# triggers a systemd-analyze verify warning and is ignored on current
# systemd. Keeps restart-loop cap effective for lightdm.service.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
# systemd creates /run/lightdm with mode 0711, owner lightdm:lightdm.
# The lightdm pkg's sysusers.d creates that user.
RuntimeDirectory=lightdm
RuntimeDirectoryMode=0711
# StateDirectory similarly ensures /var/lib/lightdm exists.  The package
# tmpfiles entry covers this but again, defensive.
StateDirectory=lightdm
StateDirectoryMode=0750
# LogsDirectory for /var/log/lightdm (LightDM writes seat logs there).
LogsDirectory=lightdm
LogsDirectoryMode=0750
# Restart on crash so a transient failure doesn't leave the user at a
# black console.  'on-failure' covers nonzero-exit and signals; it does
# NOT loop on clean exit (systemctl stop).
Restart=on-failure
RestartSec=2s
