# profile/airootfs/etc/systemd/system/lightdm.service.d/s80-preflight.conf
# S80 Agent X -- fail-fast preflight + diagnostic logging for lightdm.
#
# On S79 ISO boots to real hardware, some users report that the XFCE
# desktop never appears. The underlying cause was an ordering cycle
# that caused ai-hw-detect.service to be skipped, leaving the downstream
# `Wants=` chain unsatisfied (fixed separately in this session by removing
# `Before=basic.target` from ai-hw-detect).
#
# This drop-in adds belt-and-suspenders protection:
#
#   1. ConditionPathExists=/etc/lightdm/Xsession -- fails the unit FAST if
#      our Xsession wrapper is missing.  Without it LightDM starts, tries
#      to exec a missing wrapper, and loops silently.
#
#   2. ConditionPathExists=/usr/bin/startxfce4 -- fails fast if xfce4-session
#      isn't installed.  LightDM would otherwise start, call the wrapper,
#      which would try to exec `startxfce4`, which would ENOENT and the user
#      sees a black screen.
#
#   3. ExecStartPre=/bin/sh -c "..." -- logs prereq status to journal for
#      post-mortem diagnosis:  journalctl -u lightdm.service -b
#
# If this unit fails-to-start because of a Condition, systemd skips it
# cleanly and the ai-fbcon-fallback.service from fallback.conf triggers
# on OnFailure= -- but only on installed-to-disk (ConditionPathExists=
# !/run/archiso guards the live ISO from the fallback).
#
[Unit]
ConditionPathExists=/etc/lightdm/Xsession
ConditionPathExists=/usr/bin/startxfce4

[Service]
# Log prereq state to journal BEFORE LightDM starts. Non-fatal: all the
# real gates live in the Condition* above.  This gives operators a one-
# line summary to grep for when XFCE doesn't appear.
ExecStartPre=/bin/sh -c 'echo "s80-preflight: Xsession=$(test -x /etc/lightdm/Xsession && echo OK || echo MISSING) startxfce4=$(test -x /usr/bin/startxfce4 && echo OK || echo MISSING) xsess=$(ls /usr/share/xsessions/*.desktop 2>/dev/null | tr \"\\n\" \" \") arch_user=$(id arch 2>/dev/null | head -c120) hw_profile=$(test -f /run/ai-arch-hw-profile && echo OK || echo MISSING)"'
