# profile/airootfs/etc/systemd/system/lightdm.service.d/s81-xorg-safety.conf
# S81.5 Agent X2 -- Xorg-safety drop-in.
#
# Complements s80-preflight.conf and runtime-dir.conf.  Addresses the
# real-hardware symptom: after the Plymouth CRLF fix, boots reach the
# "Welcome to Arch Linux" getty banner on tty1, then cycle through a
# black screen back to the banner.  This is LightDM's classic "X died,
# respawn, X dies again" loop.
#
# Layered defenses here (file-level doc):
#
#   1. ConditionPathExists=/usr/share/xsessions/xfce.desktop
#      s80-preflight already conditions on /etc/lightdm/Xsession and
#      /usr/bin/startxfce4, but NOT on the xsession .desktop file.
#      lightdm.conf:18 has `autologin-session=xfce` which means LightDM
#      looks for /usr/share/xsessions/xfce.desktop specifically.  If
#      the customize_airootfs.sh fallback at line 196 failed to write
#      that file (somehow) LightDM would try the session, fail to
#      resolve it, exit 1, respawn -- silent loop.
#
#   2. StartLimitBurst=3 (was 5 in runtime-dir.conf).
#      With Restart=on-failure + RestartSec=2s, 5 restarts takes ~15 s
#      before systemd gives up and triggers OnFailure.  On real HW the
#      user has stared at a black screen for that time with no clue
#      whether to wait or reboot.  3 restarts is ~8 s -- faster drop to
#      the ai-fbcon-fallback.service emergency shell where they can
#      read /var/log/lightdm/ and fix it.
#      runtime-dir.conf's value is overridden by the MORE-RESTRICTIVE
#      value here per systemd drop-in merge rules (last-wins on scalar).
#
#   3. Environment=XDG_SESSION_TYPE=x11 is forced so LightDM's
#      session-type autoprobing can't accidentally select "wayland"
#      on a GPU with wayland support but no wayland session installed
#      (only xfce.desktop ships, no wayland .desktop).  Belt-and-
#      suspenders; normally autologin-session=xfce in lightdm.conf
#      already pins this.
#
#   4. ExecStartPre logs Xorg.0.log-like info BEFORE LightDM starts,
#      so a post-mortem `journalctl -u lightdm -b` reveals DRM devices
#      and xorg.conf.d/ state even when X itself never got far enough
#      to write its own log.
#
[Unit]
ConditionPathExists=/usr/share/xsessions/xfce.desktop
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Environment=XDG_SESSION_TYPE=x11
ExecStartPre=/bin/sh -c 'echo "s81-xorg-safety: drm=$(ls -1 /dev/dri/ 2>/dev/null | tr \"\\n\" \" \") xorg_confd=$(ls -1 /etc/X11/xorg.conf.d/ 2>/dev/null | grep -v disabled | tr \"\\n\" \" \") nouveau_loaded=$(lsmod 2>/dev/null | grep -qw nouveau && echo YES || echo no) nvidia_loaded=$(lsmod 2>/dev/null | grep -qw nvidia && echo YES || echo no) amdgpu_loaded=$(lsmod 2>/dev/null | grep -qw amdgpu && echo YES || echo no) i915_loaded=$(lsmod 2>/dev/null | grep -qw i915 && echo YES || echo no)"'
