# profile/airootfs/usr/lib/systemd/system/pe-compat-firewall.service.d/slice.conf
# pe-compat-firewall.service drop-in -- join ai-daemon.slice.
#
# The firewall is an observer-style process: it watches netlink for
# connection events and applies nft rules.  Its own CPU footprint is
# small but bursty when a game opens 40 connections in a second.
#
# Placement rationale:
#   - NOT observer.slice: the firewall is latency-sensitive (if we
#     throttle it, connections get made before rules apply).
#   - NOT pe-compat.slice: we want it to lose to the game, not compete.
#   - ai-daemon.slice: perfect fit -- medium priority, bounded memory,
#     still runs promptly when called.

[Service]
Slice=ai-daemon.slice

# Network-event-driven service -- responsive but not priority.
Nice=-2
OOMScoreAdjust=-100

# Realtime-ish IO for nft commit.  Without this, a loaded disk can
# delay rule installation by seconds.
IOSchedulingClass=best-effort
IOSchedulingPriority=2
