# profile/airootfs/etc/systemd/system/ModemManager.service.d/lazy.conf
# ModemManager only matters on laptops with WWAN / cellular modems.
# Most users don't have one; cap the probe timeout so a hung modem
# adapter doesn't stall boot.

[Service]
TimeoutStartSec=10
SuccessExitStatus=0 1
