# profile/airootfs/etc/systemd/system/plymouth-quit.service.d/timeout.conf
[Service]
# CLAUDE.md target: 10s ceiling on plymouth-quit. Default is infinity.
# Stuck plymouth-quit is the #1 cause of "boot hangs at splash" on this ISO.
TimeoutStartSec=10
TimeoutStopSec=10
SuccessExitStatus=0 1
