# profile/airootfs/etc/systemd/system/systemd-journal-flush.service.d/fast.conf
# systemd-journal-flush flushes /run/log/journal -> /var/log/journal.
# With Storage=volatile (see 50-boot-speed.conf) there is no persistent
# journal to flush -- the service is effectively a no-op but still adds
# 2-3 s of ordering dependencies on boot.  Short-circuit its timeout.

[Unit]
# Don't block shutdown on a unit that has nothing to do.
DefaultDependencies=no

[Service]
# If it ever does run, don't let it drag out.
TimeoutStartSec=5
SuccessExitStatus=0 1
