# profile/airootfs/etc/systemd/system/NetworkManager-wait-online.service.d/timeout.conf
# Short-circuit NetworkManager-wait-online on live ISO and headless machines.
#
# Default wait is 30 seconds (ExecStart=nm-online -s -q --timeout=30).
# That is 30 seconds of literal hang in QEMU/VM boots where the NIC comes
# up before NetworkManager even tries.  Our services do not require
# connectivity before multi-user.target -- they handle "no network" gracefully.
#
# Drop the timeout to 5 seconds and mark the unit non-essential.

[Unit]
# Still pulled in by NetworkManager.service but not a barrier for the whole boot
DefaultDependencies=no

[Service]
ExecStart=
ExecStart=/usr/bin/nm-online -s -q --timeout=5
TimeoutStartSec=10
# Don't fail the boot if the NIC isn't ready yet
SuccessExitStatus=0 1
