# profile/airootfs/etc/systemd/system/ai-driver-loader.service.d/10-hardening.conf
[Service]
# Session 68 hardening: bound capabilities from full-root to minimum needed for modprobe.
CapabilityBoundingSet=CAP_SYS_MODULE CAP_NET_ADMIN CAP_SYS_ADMIN
# Keep NoNewPrivileges OFF per existing design (can be overridden in main unit) but filter syscalls.
SystemCallArchitectures=native
SystemCallFilter=@system-service @module
SystemCallFilter=~@reboot @swap @clock
RestrictNamespaces=yes
LockPersonality=yes
