# profile/airootfs/etc/systemd/system/ai-driver-loader.service.d/plymouth-quiet.conf
# Suppress console noise during the Plymouth boot phase.
#
# ai-load-drivers logs every modprobe attempt to the journal via `logger`,
# but we don't want those lines stomping on the boot splash on slower
# machines where Plymouth is still drawing.  Send all stdout/stderr to
# null at the unit level (journal collection still happens via syslog).
#
# This drop-in is in [Service] (systemd-analyze verify clean).
[Service]
StandardOutput=null
StandardError=journal
