# profile/airootfs/etc/systemd/system/ai-power.service.d/10-hardening.conf
[Service]
# Session 68 hardening: drop CAP_SYS_ADMIN (root-equivalent) which isn't needed for scaling_governor writes.
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_SYS_MODULE
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @reboot
RestrictAddressFamilies=AF_UNIX AF_NETLINK
RestrictNamespaces=yes
LockPersonality=yes
