# ai-control/config/ai-control.service
[Unit]
Description=AI Control Daemon - Full System Control
Documentation=https://github.com/ai-arch-linux/ai-control-daemon
# Start after basic networking is available. Do NOT use network-online.target
# because it can hang indefinitely if NetworkManager hasn't established a
# connection (common in QEMU/first-boot). The daemon doesn't need full
# network connectivity to start — it binds to localhost:8420.
After=network.target
# Start after systemd basic targets are up
After=local-fs.target
# ai-hw-detect drops /run/ai-arch-hw-profile which main.py reads at startup
# to learn whether it's on OLD (P4) or NEW hardware.  Wants= not Requires=
# so a hw-detect failure doesn't block the daemon (it falls back to DEFAULT).
After=ai-hw-detect.service
Wants=ai-hw-detect.service
# StartLimit* live in [Unit] as of systemd 230.  Older placement in
# [Service] still works but triggers a deprecation warning in
# systemd-analyze verify.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
# Type=notify: the daemon calls sd_notify(0, "READY=1") AFTER uvicorn has
# finished its lifespan.startup AND is accepting connections on port 8420.
# Downstream units with After=ai-control.service can rely on /health being
# live the moment systemd reports us "active". Prior Type=simple reported
# active as soon as Python forked, which caused a readiness race (SSH up
# but /health 404 / connection-refused).
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/ai-control-daemon
StandardOutput=journal
StandardError=journal
# Run as root for full system access
User=root
Group=root
# Working directory
WorkingDirectory=/usr/lib/ai-control-daemon
# Environment — do NOT set DISPLAY here. The daemon auto-detects the
# session type (Wayland/X11/headless) at startup. Setting DISPLAY=:0
# when X isn't running yet causes compositor/screen modules to fail.
# XAUTHORITY is also auto-detected from /run/lightdm, ~/.Xauthority, etc.
Environment=PYTHONUNBUFFERED=1
# Directories
RuntimeDirectory=ai-control-daemon
StateDirectory=ai-control-daemon
LogsDirectory=ai-control-daemon
# Restart policy — restart on failure, but also on abnormal exit/watchdog
Restart=on-failure
RestartSec=3
# NOTE: StartLimitIntervalSec / StartLimitBurst moved to [Unit] above.
# Bumped from 30s → 90s to accommodate the cold-start path:
# Python module imports, FastAPI app construction, cortex/pattern-scanner
# lazy init, uvloop import, and uvicorn socket bind. On TCG-emulated QEMU
# this is consistently ~30-45s; 90s gives headroom for 1GB Pentium 4.
TimeoutStartSec=90
# Bounded shutdown.  main.py's signal handlers close uvicorn gracefully
# within ~10s (see asyncio.gather timeout in main._shutdown).  Give ourselves
# 15s before systemd escalates to SIGKILL so in-flight audit writes flush.
TimeoutStopSec=15
# Kill main process first; allow children to receive SIGTERM individually.
# mixed avoids the "cgroup-wide SIGKILL" path which can truncate logs.
KillMode=mixed
KillSignal=SIGTERM
# Session 36: wired. main.py runs _watchdog_heartbeat() every 30s (half
# of WatchdogSec) after READY=1, and sends STOPPING=1 on shutdown. If the
# asyncio loop stalls for >60s, systemd SIGKILLs and Restart=on-failure
# brings us back.
WatchdogSec=60
# No restrictions - AI has full access
ProtectSystem=false
ProtectHome=false
PrivateTmp=false
NoNewPrivileges=false
# Full capabilities (root default - no restrictions)

[Install]
WantedBy=multi-user.target
