# profile/airootfs/usr/lib/systemd/system/pe-objectd.service.d/hive-env.conf
# Session 38 drop-in: pin HOME to a persistent system path so registry.c
# places its hive tree under /var/lib/pe-compat/.pe-compat/registry/
# instead of /root/.pe-compat/registry/ (or /tmp on HOME-unset failures).
#
# With this, objectd_registry_init()'s canonical-mirror pass creates
# /var/lib/pe-compat/registry/{SOFTWARE,SYSTEM,HKLM,HKCU,HKCR,HKU,HKCC}
# as symlinks pointing at the real hive directories -- the path the
# QEMU smoke test [24] probes for SOFTWARE hive existence.
#
# Living under /var/lib guarantees persistence across reboots (unlike
# /root which is on a volatile overlay during live-ISO boot) and
# survives package upgrades (unlike /root/.pe-compat which XFCE session
# migration tooling occasionally purges).

[Service]
Environment=HOME=/var/lib/pe-compat
