# profile/airootfs/usr/lib/systemd/system/pe-objectd.service.d/slice.conf
# pe-objectd.service drop-in -- join pe-compat.slice.
#
# The object broker (L1) sits below PE processes and handles named
# objects, registry hives, session manager.  It's a direct dependency
# of every PE app, so it must share scheduling priority with them.
#
# We explicitly do NOT place it in trust.slice: it's NOT an authority
# root, just a convenience broker.  A hung broker blocks PE apps but
# doesn't corrupt trust state.

[Service]
Slice=pe-compat.slice

# Higher priority than random PE apps because every app talks to the
# broker on startup.  Nice=-3 gives it a small scheduling boost without
# stealing from the game itself.
Nice=-3

# Don't let the broker be killed lightly -- killing it cascades to every
# PE process losing its handles.
OOMScoreAdjust=-200
