# profile/airootfs/usr/lib/systemd/system/ai-cortex.service.d/slice.conf
# ai-cortex.service drop-in -- join ai-daemon.slice.
#
# The cortex is the decision engine (L4).  It sits above the daemon in
# the stack and mostly reads from it, so they belong in the same slice:
# contention between them is fine, but together they must lose to
# pe-compat when a game runs.

[Service]
Slice=ai-daemon.slice

# The cortex hosts the autonomy controller which chooses approvals;
# when a game is launched it's the first thing to decide "yes demote me".
# OOMScoreAdjust=-50 (less protected than ai-control since a cortex
# restart is cheaper -- it doesn't reload the LLM).
OOMScoreAdjust=-50

# No Nice here -- the cortex runs analysis on a 1 Hz tick; giving it
# interactive priority would waste cycles.
