# profile/airootfs/etc/systemd/system/systemd-logind.service.d/10-plymouth-order.conf
# AI Arch Linux — order systemd-logind after plymouth-start
#
# Motivation:
#   If the user stabs the power button while the boot splash is still the
#   first frame, logind may have been listening on /dev/input/eventN before
#   Plymouth owned the framebuffer. Shutting down then briefly unblanks the
#   console and the kernel's own init/sync messages leak through between
#   the splash and the poweroff.
#
#   Ordering logind After=plymouth-start.service guarantees the splash owns
#   the screen before the button can hand us off to plymouth-poweroff.service,
#   so the teardown path is: splash -> splash (shutdown theme) -> power off,
#   with no bare-console flash in between.
#
#   We do NOT add Requires= — if Plymouth is disabled (debug boot), logind
#   must still start. After= is a soft order, not a dependency.

[Unit]
After=plymouth-start.service
After=systemd-udev-settle.service

[Service]
# Make logind's own exit quieter during shutdown so it doesn't race the
# splash-teardown with a stderr burst.
StandardOutput=null
StandardError=journal
