# profile/airootfs/usr/lib/systemd/system/ai-hw-detect.service
[Unit]
Description=AI Arch hardware profiler (old vs new HW detection)
# Runs as part of basic.target so /run/ai-arch-hw-profile exists before
# any service that depends on it (display-manager, ai-control, firewall).
# DefaultDependencies=yes so it gets proper early-boot ordering for free.
# After=swap.target ensures zram (if enabled) is up before sysctl-tune.sh
# probes /proc/swaps for zram-aware swappiness.
After=systemd-modules-load.service local-fs.target swap.target
# S80: REMOVED `Before=basic.target`. That created an ordering cycle on
# real-hardware boot: swap.target is a sysinit.target dependency, and the
# implicit `After=sysinit.target` (DefaultDependencies=yes) combined with
# `Before=basic.target` forced this unit between sysinit and basic, but
# ai-low-ram-services/ai-irq-balance transitively depend on local-fs
# (which basic pulls). systemd broke the cycle by skipping this service
# → display-manager saw no /run/ai-arch-hw-profile → XFCE/lightdm hung.
# Services that genuinely depend on the profile (display-manager, lightdm,
# ai-control) already carry their own explicit After= on this unit.
Before=display-manager.service ai-control.service lightdm.service \
       ai-low-ram-services.service
Documentation=file:///usr/lib/ai-arch/hw-detect.sh

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/ai-arch/hw-detect.sh
# Fast: 10-second ceiling.  If detection hangs (e.g., missing lspci)
# the unit quietly falls back to DEFAULT profile and boot continues.
TimeoutStartSec=10
# Fail silently -- never block boot on a classifier hiccup.
SuccessExitStatus=0 1 2
# Lightweight sandbox (still allow writing to /etc + /sys)
ProtectSystem=no
NoNewPrivileges=yes
# Give detection CPU to run quickly on loaded systems
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=4

[Install]
WantedBy=multi-user.target graphical.target
