# profile/airootfs/etc/systemd/system/plymouth-quit-wait.service.d/timeout.conf
[Service]
# Cap plymouth-quit-wait so a stuck splash can't block the boot.
# CLAUDE.md target: 15s on plymouth-quit-wait, 10s on plymouth-quit.
# Default is infinity (waits for all VT handovers); on live ISO / headless
# QEMU / fbcon-only hosts there is no VT transition and we would hang.
TimeoutStartSec=15
TimeoutStopSec=15
# If the splash is still alive at this point, succeed anyway — it will be
# torn down by plymouth-quit which has its own 10s cap.
SuccessExitStatus=0 1
