# profile/airootfs/etc/systemd/system/plymouth-quit-wait.service.d/timeout.conf
[Service]
# Cap plymouth-quit-wait so a stuck splash can't block the boot.
# S82+D: 15s on TimeoutStartSec (wait for splash handoff), 10s on
# TimeoutStopSec (force-kill if plymouthd won't quit). Default is
# infinity (waits for all VT handovers); on live ISO / headless QEMU /
# fbcon-only hosts there is no VT transition and we would hang.
# If plymouth is stuck, force-kill at 10s after a 15s handoff wait so
# lightdm.service (After=plymouth-quit-wait.service) cannot starve.
TimeoutStartSec=15
TimeoutStopSec=10
# If the splash is still alive at this point, succeed anyway — it will be
# torn down by plymouth-quit which has its own 10s cap.
SuccessExitStatus=0 1
