# profile/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Session 38 agent-2: getty@tty1 autologin drop-in for live ISO.
#
# Why this file exists:
#  1) Archiso upstream ships `/etc/systemd/system/getty@tty1.service.d/
#     autologin.conf` that autologins ROOT on tty1 with --noclear. On a
#     graphical.target system with LightDM grabbing vt7 (our lightdm.conf
#     has minimum-vt=7), this is fine in theory — but some users observed
#     the desktop not loading on bare metal.
#  2) We override to autologin the `arch` user (not root) so the debug
#     shell matches the LightDM autologin session (both become `arch`).
#  3) Critical: we do NOT add Conflicts=graphical.target or After=
#     display-manager.service, so the getty never blocks DM progression.
#     The getty runs on tty1 INDEPENDENTLY from LightDM (which owns tty7).
#  4) --noclear preserved so boot messages remain visible for debug.
#
# If the desktop fails to load, the user can switch to tty1 with Ctrl+Alt+F1
# and get a logged-in `arch` shell to diagnose (journalctl -u lightdm, etc.)

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin arch --noissue %I $TERM
Type=idle
