# profile/airootfs/usr/lib/systemd/system/scm-daemon.service.d/slice.conf
# scm-daemon.service drop-in -- join pe-compat.slice.
#
# Windows Service Control Manager (L3).  Manages Win32 service lifecycle,
# driver hosts, dependency graph.  Shares slice with pe-compat because
# SCM is part of the Windows runtime stack; trust is a lower layer that
# scm relies on.

[Service]
Slice=pe-compat.slice

# Slight priority boost -- SCM mediates inter-service startup, any lag
# here propagates to every Windows service booting.
Nice=-2
OOMScoreAdjust=-150
