# profile/airootfs/etc/systemd/system/ai-control.service.d/group.conf
## ai-control.service drop-in: non-root accessibility of the event bus
#
# Session 69 (Agent R).  The S68 audit flagged that the ai-control daemon
# runs as User=root and creates /run/pe-compat/events.sock at 0660
# root:root, which silently locks out every non-root caller (the ISO's
# default "arch" user, any user created by ai-install-to-disk).
#
# Fix without changing the upstream unit: set SupplementaryGroups so the
# daemon's process GID list contains both "trust" and "pe-compat", then
# chgrp the socket to pe-compat from inside event_bus.py at bind time.
# With this drop-in:
#
#   * The trust group grants read ioctls on /dev/trust to the daemon's
#     supplementary GIDs, making trust-band queries work from userspace
#     helpers spawned by the daemon without a setuid shim.
#
#   * The pe-compat group on /run/pe-compat/events.sock lets members of
#     that group subscribe to the event bus (cortex, diag tools, future
#     CLI probes) without escalating to root.
#
# NOTE: we do NOT drop privileges here (still User=root).  The daemon
# still needs CAP_NET_ADMIN for nftables, CAP_SYS_PTRACE for window
# enumeration, and the trust kernel module's write-side ioctls.  What
# this buys us is that the group-0660 sockets/devices the daemon
# creates are accessible to non-root members of those groups.
[Service]
SupplementaryGroups=trust pe-compat
