# profile/airootfs/etc/systemd/system/systemd-networkd-wait-online.service.d/timeout.conf
# Short-circuit systemd-networkd-wait-online timeout (defaults to 120 s).
# This unit is ordered Before=network-online.target; if anything in the
# system pulls it in (it ships enabled in modern systemd), it can stall
# boot for up to two minutes on disconnected machines.

[Service]
ExecStart=
ExecStart=/usr/lib/systemd/systemd-networkd-wait-online --timeout=5
TimeoutStartSec=10
SuccessExitStatus=0 1
