# profile/airootfs/usr/lib/systemd/system/ai-control.service.d/slice.conf
# ai-control.service drop-in: place the daemon in ai-daemon.slice so it
# inherits the CPUWeight/MemoryHigh/IOWeight tuning defined there.
#
# This file lives in /usr/lib/systemd/system/ai-control.service.d/ so it
# ships with the ai-control-daemon package *and* the archiso profile -- on
# an installed system both locations will have the drop-in applied.
#
# Why a drop-in instead of editing ai-control.service?
#   The .service file is emitted by packages/ai-control-daemon/PKGBUILD
#   from ai-control/config/ai-control.service which is outside THIS
#   agent's scope (and rewriting it on every resource tweak is noisy).
#   Drop-ins are the standard systemd way to layer overrides.
#
# Why not systemd 254+ Slice=ai-daemon.slice inline?
#   We want the override to take effect even on systems that pacstrap'd
#   an older ai-control-daemon that predates the slice work -- drop-ins
#   work back to systemd 215.

[Service]
Slice=ai-daemon.slice

# Per-service OOM score: ai-control is not critical (a game can run
# without it) but its restart on OOM is expensive (loads LLM, rebuilds
# fastapi route table, reopens audit log).  Middle ground: -100 means
# the kernel prefers almost anything else but won't go nuts protecting
# it when observer.slice and user processes are available.
OOMScoreAdjust=-100

# Nice: slight boost.  The daemon handles websocket traffic from the
# user's frontend and users notice sub-100ms hitching there.  Nice=-5
# tips the scheduler toward giving the daemon a CPU when it wakes up.
Nice=-5

# IO priority: best-effort, priority 4 (default).  The daemon's disk
# writes are all audit/log related and non-latency-critical.
IOSchedulingClass=best-effort
IOSchedulingPriority=4
