# profile/airootfs/etc/systemd/system/NetworkManager.service.d/fast-start.conf
# NetworkManager ships with TimeoutStartSec=90s.  On a live ISO with
# no pre-configured connections there's nothing to do beyond talking
# to iwd; 15 s is ample.

[Service]
TimeoutStartSec=15
# Speed up DHCP retry loop -- default 45 s is too long for live boot.
Environment=NM_DEBUG=0
