# profile/airootfs/etc/systemd/system/systemd-udev-settle.service.d/fast.conf
# systemd-udev-settle is deprecated on modern systemd but still pulled in
# by lvm2, mdadm, and some legacy drivers.  Default timeout is 180 s;
# on a live ISO with no MD arrays or exotic storage, it either completes
# fast or never will.  Cap at 10 s.

[Service]
TimeoutStartSec=10
SuccessExitStatus=0 1
