# profile/airootfs/etc/systemd/system/bluetooth.service.d/lazy.conf
# Bluetooth: keep ENABLED for laptops but don't block boot on adapter probe.
# Some USB BT dongles take 8-12 s to enumerate; the default 90 s start
# timeout can block graphical.target on a broken adapter.
#
# OLD-profile hosts mask this unit entirely via low-ram-services.sh.

[Service]
TimeoutStartSec=10
# Always succeed -- a missing adapter is not a boot-blocker.
SuccessExitStatus=0 1
