# profile/airootfs/etc/systemd/system/plymouth-poweroff.service.d/splash.conf
[Unit]
# Ensure the splash is up BEFORE systemd starts shutting services down.
# Without this, the service races systemd-poweroff and the user sees a
# flash of "kernel initialization" text before the screen goes dark.
DefaultDependencies=no
Before=poweroff.target systemd-poweroff.service
Conflicts=plymouth-quit.service plymouth-quit-wait.service

[Service]
# Same 10s ceiling we use on plymouth-quit — a stuck splash must never
# block poweroff.
TimeoutStartSec=10
TimeoutStopSec=10
SuccessExitStatus=0 1

[Install]
WantedBy=poweroff.target
