# profile/airootfs/etc/systemd/system/serial-getty@ttyS0.service.d/autologin.conf
# Session 38 agent-2: serial-getty autologin for QEMU smoke test.
#
# The scripts/test-qemu.sh harness relies on SSH (port 2222) for the real
# tests, but having an autologin serial console is invaluable for manual
# `qemu -serial stdio` debugging when SSH itself is broken.
#
# This drop-in is ONLY activated when serial-getty@ttyS0.service is started
# (i.e., when the kernel sees a ttyS0 device — QEMU -serial, real serial
# hardware). On headless bare-metal it is inert.

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --keep-baud --autologin arch 115200,57600,38400,9600 %I $TERM
Type=idle
