# profile/airootfs/etc/systemd/system/pe-compat-firewall.service.d/wants-slice.conf
# Session 44: the firewall module enforces cgroup-based policy under
# /sys/fs/cgroup/pe-compat.slice.  Without a hard Wants= the slice unit
# is not activated until the first PE scope is spawned under it, so
# firewall's ensure_subtree_controllers() logs ENOENT on every boot
# before the first game launch.  Pulling in pe-compat.slice here makes
# the directory exist by the time nft_manager opens its netlink socket.
[Unit]
Wants=pe-compat.slice
After=pe-compat.slice
